# sv/unclutter/run
#!/bin/sh

exec unclutter -jitter 10 -idle 2
